# integration-test/src/travis/system/certhub-certbot-run@.service.d/hook-lexicon-auth.conf.in
[Service]
Environment="CERTHUB_LEXICON_GLOBAL_ARGS=--delegated ci.certhub.io"
Environment="CERTHUB_LEXICON_PROVIDER=digitalocean"
Environment="LEXICON_DIGITALOCEAN_TOKEN=${LEXICON_DIGITALOCEAN_TOKEN}"
