# lib/systemd/certhub-cert-reload@.service.d/reload.conf
[Service]

# CERTHUB_CERT_SERVICES_RELOAD
# Path to a file containing the services to reload.
#
Environment=CERTHUB_CERT_SERVICES_RELOAD=/home/certhub/config/%i.services-reload.txt
