# lib/systemd/system/certhub-dehydrated-run@.path
[Unit]
Description=Trigger certificate renewal via certhub/dehydrated when expiry status changed
Documentation=man:dehydrated(1)

[Install]
WantedBy=paths.target
