# lib/systemd/system/certhub-certbot-run@.path
[Unit]
Description=Trigger certificate renewal via certhub/certbot when expiry status changed
Documentation=man:certbot(1)

[Install]
WantedBy=paths.target
