# integration-test/src/travis/system/certhub-lego-run@.service.d/lego-challenge.conf.in
[Service]
Environment="LEGO_CHALLENGE_ARGS=--dns digitalocean"
Environment="DO_AUTH_TOKEN=${LEXICON_DIGITALOCEAN_TOKEN}"
