# lib/systemd/certhub-cert-expiry@.path
[Unit]
Description=Trigger check for certificates which are about to expire when repo changed

[Install]
WantedBy=paths.target
