# lib/systemd/certhub-cert-expiry@.path
[Unit]
Description=Trigger check for certificates which are about to expire when repo changed

[Path]
PathChanged=/home/certhub/certs.git/refs/heads/master

[Install]
WantedBy=paths.target
