# lib/systemd/certhub-cert-reload@.path
[Unit]
Description=Trigger services reload when certificate changed

[Install]
WantedBy=paths.target
