# lib/systemd/certhub-cert-send@.path
[Unit]
Description=Send certificates to predefined destinations when certificate changed

[Install]
WantedBy=paths.target
