# lib/systemd/certhub-lego-run@.path
[Unit]
Description=Trigger certificate renewal via certhub/lego when expiry status changed

[Install]
WantedBy=paths.target
