# lib/systemd/certhub-cert-export@.path
[Unit]
Description=Trigger certificate export when repo changed

[Install]
WantedBy=paths.target
