# lib/systemd/certhub-cert-reload@.path
[Unit]
Description=Trigger services reload when certificate changed

[Path]
PathChanged=/etc/ssl/certhub/%i/fullchain.pem

[Install]
WantedBy=paths.target
