# lib/systemd/certhub-mirror-push@.path
[Unit]
Description=Push certhub repo to another host when HEAD changed
Documentation=man:git-push(1)

[Path]
PathChanged=/home/certhub/certs.git/HEAD

[Install]
WantedBy=paths.target
