# lib/systemd/certhub-repo-push@.service.d/remote.conf
[Service]

# CERTHUB_REPO_REMOTE
# URL of the remote repository.
#
Environment=CERTHUB_REPO_REMOTE=%I

# CERTHUB_REPO_PUSH_ARGS
# Arguments to the git push command.
# See man:git-push(1)
#
Environment=CERTHUB_REPO_PUSH_ARGS=--mirror

# CERTHUB_REPO_PUSH_REFSPEC
# Refspec for the git push command.
# See man:git-push(1)
#
Environment=CERTHUB_REPO_PUSH_REFSPEC=

