# lib/systemd/certhub-cert-export@.service.d/export.conf
[Service]
# CERTHUB_EXPORT_SRC
# File / directory inside the repository which should be exported.
#
Environment=CERTHUB_EXPORT_SRC={WORKDIR}/%i.fullchain.pem

# CERTHUB_EXPORT_DEST
# File / directory where the certificate should be placed.
#
Environment=CERTHUB_EXPORT_DEST=/etc/ssl/certhub/%i.fullchain.pem

# CERTHUB_RSYNC_ARGS
# Arguments for rsync.
#
Environment="CERTHUB_RSYNC_ARGS=\
    --checksum \
    --delete \
    --devices \
    --links \
    --perms \
    --recursive \
    --specials \
    --verbose"
