# lib/systemd/certhub-cert-expiry@.service.d/expiry.conf
[Service]

# CERTHUB_CERT_EXPIRY_TTL:
# See man:certhub-cert-expiry(1)
#
Environment=CERTHUB_CERT_EXPIRY_TTL=2592000

# CERTHUB_CERT_EXPIRY_MESSAGE:
# Message written to the status file if certificate is about to expire.
#
Environment="CERTHUB_CERT_EXPIRY_MESSAGE=Certificate will expire within 30 days"

# CERTHUB_MESSAGE_SUBJECT:
# First line of the status file. By default this is generated automatically.
# See man:certhub-message-format(1)
#
Environment="CERTHUB_MESSAGE_SUBJECT=[Certhub] Stale %i"
