# lib/systemd/certhub-dehydrated-run@.service.d/dehydrated.conf
[Service]
Type=oneshot

# DEHYDRATED_ARGS:
# Additional Arguments for dehydrated --signcsr (see man:dehydhrated(1))
#
Environment=DEHYDRATED_ARGS=

# DEHYDRATED_CONFIG:
# Path to a dehydrated configuration file.
#
Environment=DEHYDRATED_CONFIG=/home/certhub/config/%i.dehydrated
